# scripts/cron/registration_cleanup.cron
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
CRON_TZ=UTC

0 * * * * /usr/local/bin/python -m scripts.registration_cleanup
